# deploy/backup/tronus-backup.service
[Unit]
Description=Tronus Postgres nightly backup
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
ExecStart=/opt/tronus/deploy/backup/backup-postgres.sh
StandardOutput=journal
StandardError=journal
